# root/etc/services.d/apt-cacher-ng/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid abc apt-cacher-ng -c /config/conf ForeGround=1